# etc/cron.weekly/cron.weekly-restart-phigita.sh
#!/bin/sh

/web/bin/restart-phigita.sh


